# config/services/beat.conf
[unix_http_server]
file=/dev/shm/supervisor.sock   ; (the path to the socket file)

[supervisord]
logfile=/opt/logs/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=1024MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=1           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/dev/shm/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=100                  ; (min. avail startup file descriptors;default 1024)
minprocs=500                 ; (min. avail process descriptors;default 200)
[rpcinterface:supervisor]
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///dev/shm/supervisor.sock ; use a unix:// URL  for a unix socket

[supervisord]
nodaemon=true

[program:beat_trade_smart]
environment=
         C_FORCE_ROOT=true
command=bash -c "/usr/local/bin/celery -A celery_app beat -l info"
directory=/usr/local/trade_smart/trade_smart_backend/trade_smart_backend/
autostart=true
autorestart=true
startsecs=10                ; number of secs prog must stay running (def. 10)
startretries=3              ; max # of serial start failures (default 3)
exitcodes=1
stopsignal=TERM             ; signal used to kill process (default TERM)
stopwaitsecs=600             ; max num secs to wait before SIGKILL (default 10)
user=root               ; setupid to this UNIX account to run the program
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0
