# config/services/web.conf
[unix_http_server]
file=/dev/shm/supervisor.sock   ; (the path to the socket file)

[supervisord]
logfile=/opt/logs/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=1024MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=1           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/dev/shm/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=100                  ; (min. avail startup file descriptors;default 1024)
minprocs=500                 ; (min. avail process descriptors;default 200)
[rpcinterface:supervisor]
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///dev/shm/supervisor.sock ; use a unix:// URL  for a unix socket

[supervisord]
nodaemon=true

[program:web]
environment=
         C_FORCE_ROOT=true
command=bash -c "/usr/local/python3.10/bin/uwsgi --ini /etc/trade_smart_uwsgi.ini"
autostart=true
autorestart=true
logfile_maxbytes=1024MB        ; max # logfile bytes b4 rotation (default 50MB)
logfile_backups=3
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0
exitcodes=1

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0
